# node.conf
start on startup
exec /usr/bin/nodejs /root/node-red-0.4.0/red.js
